# deploy/tackle-forger-backup.timer
[Unit]
Description=Daily Tackle Forger workspace backup

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=15m
Persistent=true
Unit=tackle-forger-backup.service

[Install]
WantedBy=timers.target
